# streamlit.ini
[program:streamlit]
command=/home/exouser/mambaforge/bin/streamlit run /home/exouser/pchem/apps/pchemapp.py --server.enableCORS=false --server.enableXsrfProtection=false
priority=1
autorestart=true
